# xdg/.config/autostart/dbus-user-target@sleep.desktop
[Desktop Entry]
Type=Application
Name=Dbus user target for sleep
Exec=bash -c '~/.local/bin/dbus-user-target@sleep'
X-GNOME-Autostart-enabled=true
