# xdg/.local/share/systemd/user/sleep.target
[Unit]
Description=Sleep
StopWhenUnneeded=yes
